# trade-manager/systemd/openclaw-trade-manager.service
# Reference template — the actual service file is generated at install time
# by `trade-manager setup-service` with correct paths baked in.
#
# To install manually:  trade-manager setup-service
# To remove:            trade-manager remove-service
# Logs:                 journalctl --user -u openclaw-trade-manager -f

[Unit]
Description=OpenClaw Trade Manager
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/env node /path/to/dist/cli.js start
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=HOME=%h
WorkingDirectory=%h
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
